# systemd/tokmeter.service
[Unit]
Description=tokmeter local LLM usage proxy
After=network.target

[Service]
Type=simple
ExecStart=%h/dev/tokmeter/.venv/bin/tokmeter serve
Environment=TOKMETER_PORT=8079
Environment=TOKMETER_UPSTREAM=http://127.0.0.1:8080
# Multiple backends: replace the two lines above with PORT=UPSTREAM pairs, e.g.
# Environment=TOKMETER_LISTENERS=8079=http://127.0.0.1:8080,8000=http://127.0.0.1:8010
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
